# private_dot_config/autostart/private_startscript.sh.desktop
[Desktop Entry]
Exec=/home/rcard/Documents/startscript.sh
Icon=dialog-scripts
Name=startscript.sh
Path=
Type=Application
X-KDE-AutostartScript=true
