# private_dot_config/autostart/private_watcher.sh.desktop
[Desktop Entry]
Exec=/home/rcard/.konsole/watcher.sh
Icon=dialog-scripts
Name=watcher.sh
Path=
Type=Application
X-KDE-AutostartScript=true
